# .devops/queue-worker/supervisord.conf

[program:horizon]
process_name=%(program_name)s
command=/usr/bin/php /app/artisan horizon
autostart=true
autorestart=true
#user=ec2-user
redirect_stderr=true
stdout_logfile=/app/horizon.log
